# hw_less8/mytimer.service
[Unit]
Description=Runs unit.sh every 30 seconds
[Service]
Type=oneshot
ExecStart=/opt/unit.sh